# roles/editor/templates/editor.service.j2
# /lib/systemd/system/editor.service
[Unit]
Description=Galaxy Workflow Editor

[Service]
User=root
Type=simple
Restart=on-failure
ExecStart={{ galaxy_directory }}/run.sh